# confluent/upstart/confluent-rest.conf
description "Confluent REST"
author "Joseph <joseph.jeganathan@bedegaming.com>"
 
start on runlevel [2345]
stop on runlevel [016]

respawn

setuid vagrant

script
  kafka-rest-start
end script
